# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f
stdout_logfile=/var/log/cron.log
stderr_logfile=/var/log/cron.log

[program:django]
command=sh -c "python manage.py migrate && python manage.py runserver 0.0.0.0:3000"
stdout_logfile=/var/log/django.log
stderr_logfile=/var/log/django.log